# pi/pizero-usb-gadget/netctl/netctl@usbpizero.service.d/profile.conf
[Unit]
Description=pizero g_ether gadget
BindsTo=sys-subsystem-net-devices-usb0.device
After=sys-subsystem-net-devices-usb0.device
